# skel/.local/share/systemd/user/eating.target
[Unit]
Description=Eating...
Conflicts=thinking.target
ConditionPathIsDirectory=/home/share/dining-room
Wants=eating.timer sticky-fork.service
