# skel/.local/share/systemd/user/eating.target
[Unit]
Description=Eating...
Conflicts=thinking.target
Wants=eating.timer
