# skel/.local/share/systemd/user/fork-check.path
[Unit]
Description=Watch Forks
PartOf=hungry.target

[Path]
PathModified=/home/share/dining-room/forks
