# skel/.local/share/systemd/user/hungry.target
[Unit]
Description=Hungry...
ConditionPathIsDirectory=/home/share/dining-room
Wants=fork-check.path fork-request.service
