# skel/.local/share/systemd/user/select-seat.service
[Unit]
Description=philosopher select seat service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/select-seat
SyslogFacility=local0

[Install]
WantedBy=default.target
