# skel/.local/share/systemd/user/select-seat.path
[Unit]
Description="check if the dining room is open"

[Path]
PathModified=/home/share/dining-room

[Install]
WantedBy=default.target
