# skel/.local/share/systemd/user/thinking.target
[Unit]
Description=Thinking...
Conflicts=eating.target
ConditionPathExists=/home/share/dining-room
Wants=thinking.timer
