# mkosi.extra/etc/update-motd.d/00-header
#!/bin/sh

printf "Welcome to this Dining Philosophers simulation\n"
