# mkosi.extra/etc/systemd/system/inetutils-telnetd@.service.d/no-hostinfo.conf
[Service]
ExecStart=
ExecStart=-/usr/sbin/telnetd --no-hostinfo
