# mkosi.extra/etc/systemd/system/inetutils-telnetd.socket.d/limits.conf
[Socket]
MaxConnectionsPerSource=3
MaxConnections=16
